# crontab
00 18 * * * python /app/src/daily.py
# Don't remove the empty line at the end of this file. It is required to run the cron job
