# crontab/crontab-job
17 21 * * * /app/crontab/run_daily.sh > /app/log/crontab.log
# Don't remove the empty line at the end of this file. It is required to run the cron job
